# docker/supervisor-app.conf
[program:app-uwsgi]
command=uwsgi --ini /app/docker/uwsgi.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
